# systemd/pia-port-forward.service
[Unit]
Description=PIA port forwarding manager
After=network-online.target pia-vpn.service
Wants=network-online.target

[Service]
Type=simple
RemainAfterExit=no
WorkingDirectory=/usr/local/bin/manual-connections
EnvironmentFile=-/var/lib/pia/region.txt

# Wait for VPN to be ready
ExecStartPre=/bin/bash -c 'timeout 180 bash -c '"'"'while ! ([ -f /var/lib/pia/region.txt ] && grep -q "^gateway=" /var/lib/pia/region.txt && ip link show pia >/dev/null 2>&1 && ip addr show pia | grep -q "inet "); do sleep 1; done'"'"''

# Run port forwarding
ExecStart=/usr/local/bin/pia-port-forward-wrapper.sh

# Update firewall
ExecStartPost=/bin/bash -c 'sleep 1 && /usr/local/bin/pia-firewall-update-wrapper.sh'

# Simple restart policy
Restart=on-failure
RestartSec=10

User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
