# systemd/pia-port-forward.service
[Unit]
Description=PIA port forwarding manager
After=network-online.target pia-vpn.service
Wants=network-online.target
BindsTo=pia-vpn.service

[Service]
Type=simple
ExecStart=/usr/local/bin/pia-port-forward-check.sh
Restart=always
RestartSec=15
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
