# systemd/pia-port-forward.service
[Unit]
Description=PIA port forwarding manager
After=network-online.target pia-vpn.service
Wants=network-online.target
Requires=pia-vpn.service

[Service]
Type=simple
RemainAfterExit=no
WorkingDirectory=/usr/local/bin/manual-connections

# Read the region file to get gateway and hostname
EnvironmentFile=-/var/lib/pia/region.txt

# Wait for VPN to be ready
ExecStartPre=/bin/bash -c 'timeout 120 bash -c '"'"'while ! ([ -f /var/lib/pia/region.txt ] && grep -q "^gateway=" /var/lib/pia/region.txt && ip link show pia >/dev/null 2>&1 && ip addr show pia | grep -q "inet "); do sleep 1; done'"'"''

# Run the port forwarding wrapper (handles PIA_PF check internally)
ExecStart=/usr/local/bin/pia-port-forward-wrapper.sh

# Update firewall when port changes (runs in background)
# This runs after ExecStart succeeds and monitors for port file changes
ExecStartPost=/bin/bash -c 'sleep 1 && /usr/local/bin/pia-firewall-update-wrapper.sh'

# Restart policy
Restart=on-failure
RestartSec=15
StartLimitInterval=300
StartLimitBurst=3

User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
